# provisionning/templates/glances.service
[Unit]
Description = Glances web server
After = network.target
 
[Service]
ExecStart = /usr/bin/glances -w
 
[Install]
WantedBy=multi-user.target