# provisionning/templates/socket_activated/systemctl/proxy.socket.j2
[Socket]
ListenStream={{ public_port }}

[Install]
WantedBy=sockets.target