# provisionning/templates/socket_activated/systemctl/proxy.service.j2
[Unit]
BindTo={{ name }}-docker.service
After={{ name }}-docker.service

[Service]

ExecStart=/lib/systemd/systemd-socket-proxyd {{ service_ip }}:{{ service_port }}
